# deploy/systemd/pmc-llm-bot.service
[Unit]
Description=Prediction Market Trading Bot — LLM forecaster
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pmc
Group=pmc
WorkingDirectory=/opt/pmc
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/pmc/.env
ExecStart=/opt/pmc/venv/bin/python /opt/pmc/cli.py run --paper --platforms both --interval 60
Restart=on-failure
RestartSec=30
# If the bot crashes 5 times in 2 minutes, stop — probably something fundamentally broken
StartLimitIntervalSec=120
StartLimitBurst=5
StandardOutput=append:/var/log/pmc/llm-bot.log
StandardError=append:/var/log/pmc/llm-bot.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pmc /var/log/pmc
PrivateTmp=true

[Install]
WantedBy=multi-user.target
